# systemd/investright-symbols.service
[Unit]
Description=InvestRight monthly symbol-universe rebuild
Documentation=https://github.com/arkapchaudhuri-investright/investright
# Regenerates static/symbols.json (the autocomplete search universe) from the
# free listing sources — the Nasdaq Trader directory + NSE's equity list. Needs
# network. Writes atomically, so a source outage leaves the previous file intact.

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/investright
ExecStart=/opt/investright/.venv/bin/python build_symbols.py

# A silently stale search list is worse than a loud failure.
OnFailure=investright-alert@%n.service
